# system/logfalcon-boot-led.service
# Heartbeat LED during boot — lets the pilot know the Pi is alive.
# Automatically stopped when logfalcon-web.service starts.
# Install to: /etc/systemd/system/logfalcon-boot-led.service

[Unit]
Description=LogFalcon Boot LED Heartbeat
DefaultDependencies=no
Before=logfalcon-web.service
Conflicts=logfalcon-web.service

[Service]
Type=simple
ExecStart=/opt/logfalcon/boot-led.sh
Restart=no

[Install]
WantedBy=sysinit.target
